# integration_test/third_party_apps_data/applications/aerospike/debian_ubuntu/install
set -e

# /artifacts/aerospike-server-community/<aerospike_version>/aerospike-server-community-<aerospike_version>-<ubuntu_version>.tgz
curl -LO aerospike.tgz https://download.aerospike.com/artifacts/aerospike-server-community/5.7.0.17/aerospike-server-community-5.7.0.17-ubuntu18.04.tgz
tar -xzf aerospike-server-community-5.7.0.17-ubuntu18.04.tgz
(cd aerospike-server-community-5.7.0.17-el7 && sudo ./asinstall)
sudo systemctl enable --now aerospike-server
