# integration_test/third_party_apps_data/applications/aerospike/debian_ubuntu/install
set -e

sudo apt-get update
sudo apt-get install -y wget

source /etc/os-release
case $ID in
    debian)
        case "$VERSION_ID" in
            9*) debian_version="debian9";;
            10*|11*) debian_version="debian10";;
            *) echo "Found version: $VERSION_ID, supported versions are 9/10/11";;
        esac
        os_version=$debian_version
        ;;
    ubuntu)
        case "$VERSION_ID" in
            18*) ubuntu_version="ubuntu18.04";;
            20*) ubuntu_version="ubuntu20.04";;
            *) echo "Found version: $VERSION_ID, supported versions are 18/20"
        esac
        os_version=$ubuntu_version
        ;;
esac

# download, extract, install, and start aerospike community edition
wget -O aerospike.tgz "https://download.aerospike.com/artifacts/aerospike-server-community/5.7.0.17/aerospike-server-community-5.7.0.17-$os_version.tgz"
tar -xvf aerospike.tgz
cd aerospike-server-community-5.7.0.17-*/
sudo ./asinstall

sudo systemctl enable aerospike
sudo systemctl start aerospike
