# integration_test/third_party_apps_data/applications/elasticsearch/debian_ubuntu/install
set -e

sudo apt-get update
sudo apt-get install -y apt-transport-https wget default-jre

wget -qO - https://artifacts.elastic.co/GPG-KEY-elasticsearch | sudo apt-key add -
echo "deb https://artifacts.elastic.co/packages/7.x/apt stable main" | \
    sudo tee /etc/apt/sources.list.d/elastic-7.x.list

sudo apt-get update
sudo apt-get install -y elasticsearch

sudo systemctl daemon-reload
sudo systemctl enable elasticsearch

sudo systemctl restart elasticsearch
