# integration_test/third_party_apps_test/applications/mongodb/debian_ubuntu/install
set -e

source /etc/os-release
sudo apt-get update
sudo apt-get install -y gnupg wget curl

oscode="$(lsb_release -sc)"
osname="$(lsb_release -si)"

# TODO: Use Focal packages on Debian 11 ARM for now.
# Focal is the newest platform that's older than Debian 11,
# meaning it shouldn't run into issues with requiring newer
# dependencies than what Debian 11 offers.
# Remove once we add support for an Ubuntu ARM platform.
if [[ "${osname}" == Debian && "${oscode}" == bullseye && "$(uname -m)" == aarch64 ]]; then
  oscode="focal"
  osname="Ubuntu"
fi

MONGO_VERSION=6.0.7
case $osname in 
  Ubuntu)
    if [[ "${VERSION_ID}" =~ ^(24) ]]; then
      MONGO_VERSION=8.0.0
      curl -fsSL https://www.mongodb.org/static/pgp/server-8.0.asc | \
        sudo gpg -o /usr/share/keyrings/mongodb-server-8.0.gpg \
        --dearmor
      echo "deb [ arch=amd64,arm64 signed-by=/usr/share/keyrings/mongodb-server-8.0.gpg ] https://repo.mongodb.org/apt/ubuntu $oscode/mongodb-org/8.0 multiverse" | sudo tee /etc/apt/sources.list.d/mongodb-org-8.0.list
    else
      wget --no-verbose --output-document=/etc/apt/trusted.gpg.d/mongodb.asc https://www.mongodb.org/static/pgp/server-6.0.asc
      echo "deb [ arch=amd64,arm64 ] https://repo.mongodb.org/apt/ubuntu $oscode/mongodb-org/6.0 multiverse" | sudo tee /etc/apt/sources.list.d/mongodb-org-6.0.list
    fi
    ;;
  Debian)
    wget --no-verbose --output-document=/etc/apt/trusted.gpg.d/mongodb.asc https://www.mongodb.org/static/pgp/server-6.0.asc
    echo "deb http://repo.mongodb.org/apt/debian $oscode/mongodb-org/6.0 main" | sudo tee /etc/apt/sources.list.d/mongodb-org-6.0.list
    ;;
  *)
    echo -n "unknown os release"
    exit 1
    ;;
esac

sudo apt-get update
sudo apt-get install -y \
    mongodb-org=${MONGO_VERSION} \
    mongodb-org-database=${MONGO_VERSION} \
    mongodb-org-server=${MONGO_VERSION} \
    mongodb-org-shell=${MONGO_VERSION} \
    mongodb-org-mongos=${MONGO_VERSION} \
    mongodb-org-tools=${MONGO_VERSION}

systemctl enable mongod

sudo systemctl restart mongod