# integration_test/third_party_apps_test/applications/mongodb3.6/debian_ubuntu/install
set -e

sudo apt-get update
sudo apt-get install -y gnupg wget
wget --no-verbose --output-document=/etc/apt/trusted.gpg.d/mongodb.asc https://www.mongodb.org/static/pgp/server-3.6.asc

oscode="$(lsb_release -sc)"
osname="$(lsb_release -si)"

case $osname in
  Ubuntu)
    echo "deb [ arch=amd64,arm64 ] https://repo.mongodb.org/apt/ubuntu $oscode/mongodb-org/3.6 multiverse" | sudo tee /etc/apt/sources.list.d/mongodb-org-3.6.list
    ;;
  Debian)
    echo "deb http://repo.mongodb.org/apt/debian $oscode/mongodb-org/3.6 main" | sudo tee /etc/apt/sources.list.d/mongodb-org-3.6.list
    ;;
  *)
    echo -n "unknown os release"
    exit 1
    ;;
esac

sudo apt-get update
sudo apt-get install -y \
    mongodb-org=3.6.23 \
    mongodb-org-server=3.6.23 \
    mongodb-org-shell=3.6.23 \
    mongodb-org-mongos=3.6.23 \
    mongodb-org-tools=3.6.23

systemctl enable mongod

sudo systemctl restart mongod
