# integration_test/third_party_apps_data/applications/rabbitmq/debian_ubuntu/install
set -e

sudo apt-get update
sudo apt-get install -y \
    curl \
    gnupg \
    apt-transport-https \
    debian-keyring \
    debian-archive-keyring 
    

curl -s https://packages.erlang-solutions.com/ubuntu/erlang_solutions.asc | \
    sudo apt-key add -

source /etc/os-release
case $VERSION_ID in
  # ubuntu
  20.04)
    echo "deb https://packages.erlang-solutions.com/ubuntu focal contrib" | sudo tee /etc/apt/sources.list.d/rabbitmq.list
    ;;
  # debian
  10)
    echo "deb https://packages.erlang-solutions.com/ubuntu bionic contrib" | sudo tee /etc/apt/sources.list.d/rabbitmq.list
    ;;
  11)
    echo "deb https://packages.erlang-solutions.com/ubuntu focal contrib" | sudo tee /etc/apt/sources.list.d/rabbitmq.list
    ;;
  *)
    echo -n "unknown version"
    exit 1
    ;;
esac

sudo apt-get update
sudo apt-get install -y erlang

curl -s \
    https://packagecloud.io/install/repositories/rabbitmq/rabbitmq-server/script.deb.sh | \
    sudo bash

sudo apt-get update
sudo apt-get install -y rabbitmq-server

sudo systemctl daemon-reload
sudo systemctl enable rabbitmq-server
sudo systemctl restart rabbitmq-server
sudo rabbitmq-plugins enable rabbitmq_management
sudo systemctl restart rabbitmq-server

sudo rabbitmqctl add_vhost dev
sudo rabbitmqctl add_user "admin" "admin"
sudo rabbitmqctl set_user_tags "admin" administrator
sudo rabbitmqctl set_permissions -p "dev" "admin" ".*" ".*" ".*"

endpoint="http://127.0.0.1:15672"

curl -i -u admin:admin \
    -H "content-type:application/json" \
    -X PUT "${endpoint}/api/exchanges/dev/webex" \
    -d'{"type":"direct","auto_delete":false,"durable":true,"internal":false,"arguments":{}}'

curl -i -u admin:admin \
    -H "content-type:application/json" \
    -X PUT "${endpoint}/api/queues/dev/webq1" \
    -d'{"auto_delete":false,"durable":true,"arguments":{}}'

curl -i -u admin:admin \
    -H "content-type:application/json" \
    -X POST "${endpoint}/api/bindings/dev/e/webex/q/webq1" \
    -d'{"routing_key":"webq1","arguments":{}}'

sudo rabbitmq-plugins enable rabbitmq_management
sudo rabbitmqctl add_user "otelu" "otelp"
sudo rabbitmqctl set_user_tags "otelu" monitoring
sudo rabbitmqctl set_permissions -p "dev" "otelu" "" "" ".*"
