# integration_test/third_party_apps_data/applications/couchdb/debian_ubuntu/install
set -e

sudo apt-get update
sudo apt install -y curl apt-transport-https gnupg

curl https://couchdb.apache.org/repo/keys.asc | gpg --dearmor | sudo tee /usr/share/keyrings/couchdb-archive-keyring.gpg >/dev/null 2>&1
source /etc/os-release
echo "deb [signed-by=/usr/share/keyrings/couchdb-archive-keyring.gpg] https://apache.jfrog.io/artifactory/couchdb-deb/ ${VERSION_CODENAME} main" \
    | sudo tee /etc/apt/sources.list.d/couchdb.list >/dev/null

sudo apt-get update
sudo DEBIAN_FRONTEND=noninteractive apt-get install -y couchdb

cat << EOF > local.ini
[couchdb]
[couch_peruser]
[chttpd]
[httpd]
[chttpd_auth]
[ssl]
[vhosts]
[admins]
admin = otelp
EOF
sudo mv local.ini /opt/couchdb/etc/local.ini
sudo chown couchdb:couchdb /opt/couchdb/etc/local.ini
sudo chmod 0755 /opt/couchdb/etc/local.ini

sudo systemctl enable couchdb
