# integration_test/third_party_apps_data/applications/flink/debian_ubuntu/install
set -e

sudo apt-get update
sudo apt-get install -y curl default-jre

curl -L -o flink.tgz https://dlcdn.apache.org/flink/flink-1.14.4/flink-1.14.4-bin-scala_2.11.tgz

sudo mkdir /opt/flink
sudo tar -xzf flink.tgz -C /opt/flink --strip-components 1

cat <<EOF | sudo tee /etc/systemd/system/flink.service
[Unit]
Description=Apache Flink
           
[Service]
Type=forking
User=root
Group=root
ExecStart=/opt/flink/bin/start-cluster.sh
ExecStop=/opt/flink/bin/stop-cluster.sh
Restart=always
           
[Install]
WantedBy=default.target
EOF

sudo systemctl enable flink
