# integration_test/third_party_apps_test/applications/dcgmv1/debian_ubuntu/install
set -e
source /etc/os-release
DISTRIBUTION=$(echo $ID$VERSION_ID | sed -e 's/\.//g')

# Fix for broken Bullseye Backports (prevents apt update failure)
if [[ "$DISTRIBUTION" == "debian11" ]]; then 
    if grep -q "bullseye-backports" /etc/apt/sources.list /etc/apt/sources.list.d/* 2>/dev/null; then
        sudo sed -i '/bullseye-backports/s/^/#/' /etc/apt/sources.list
        sudo sed -i '/bullseye-backports/s/^/#/' /etc/apt/sources.list.d/*.list
    fi
fi

if ! dpkg -s cuda-keyring >/dev/null 2>&1; then
    filename="cuda-keyring_1.1-1_all.deb"
    url="https://developer.download.nvidia.com/compute/cuda/repos/${DISTRIBUTION}/x86_64/${filename}"
    
    wget --no-verbose "$url"
    sudo dpkg -i "$filename"
fi

# Install DCGM
sudo apt-get update
sudo apt-get install -y datacenter-gpu-manager
sudo systemctl --now enable nvidia-dcgm

# check if DCGM service is running
# This command is only used for informational/debugging output.
dcgmi discovery --list || true
