# integration_test/third_party_apps_data/applications/elasticsearch/centos_rhel/install
set -e

sudo yum install -y java-1.8.0-openjdk

sudo rpm --import https://artifacts.elastic.co/GPG-KEY-elasticsearch

cat << EOF > elasticsearch.repo
[elasticsearch]
name=Elasticsearch repository for 7.x packages
baseurl=https://artifacts.elastic.co/packages/7.x/yum
gpgcheck=1
gpgkey=https://artifacts.elastic.co/GPG-KEY-elasticsearch
enabled=0
autorefresh=1
type=rpm-md
EOF

sudo mv elasticsearch.repo /etc/yum.repos.d/elasticsearch.repo
sudo chown root:root /etc/yum.repos.d/elasticsearch.repo
sudo chmod 0644 /etc/yum.repos.d/elasticsearch.repo

sudo yum install -y --enablerepo=elasticsearch elasticsearch 

sudo systemctl daemon-reload
sudo systemctl enable elasticsearch

sudo systemctl restart elasticsearch
