# integration_test/third_party_apps_data/applications/mongodb/debian_ubuntu/install
set -e

sudo apt-get update
sudo apt-get install -y gnupg wget
wget --no-verbose --output-document=/etc/apt/trusted.gpg.d/mongodb.asc https://www.mongodb.org/static/pgp/server-6.0.asc

oscode="$(lsb_release -sc)"
osname="$(lsb_release -si)"

# TODO: Use Focal packages on Debian 11 ARM for now.
# Focal is the newest platform that's older than Debian 11,
# meaning it shouldn't run into issues with requiring newer
# dependencies than what Debian 11 offers.
# Remove once we add support for an Ubuntu ARM platform.
if [[ "${osname}" == Debian && "${oscode}" == bullseye && "$(uname -m)" == aarch64 ]]; then
  oscode="focal"
  osname="Ubuntu"
fi

case $osname in 
  Ubuntu)
    echo "deb [ arch=amd64,arm64 ] https://repo.mongodb.org/apt/ubuntu $oscode/mongodb-org/6.0 multiverse" | sudo tee /etc/apt/sources.list.d/mongodb-org-6.0.list
    ;;
  Debian)
    echo "deb http://repo.mongodb.org/apt/debian $oscode/mongodb-org/6.0 main" | sudo tee /etc/apt/sources.list.d/mongodb-org-6.0.list
    ;;
  *)
    echo -n "unknown os release"
    exit 1
    ;;
esac

sudo apt-get update
sudo apt-get install -y \
    mongodb-org=6.0.7 \
    mongodb-org-database=6.0.7 \
    mongodb-org-server=6.0.7 \
    mongodb-org-shell=6.0.7 \
    mongodb-org-mongos=6.0.7 \
    mongodb-org-tools=6.0.7

systemctl enable mongod

sudo systemctl restart mongod
